# deb/start-stop/DEBIAN/postinst
#!/bin/bash

set -e

systemctl daemon-reload

service start-stop start

systemctl enable start-stop.service

exit 0
