# deb/start-stop/DEBIAN/prerm
#!/bin/bash

set -e

service start-stop stop
systemctl disable start-stop.service

exit 0
